# templates/etc/profile.d/mssql.sh.j2
{{ ansible_managed | comment }}
PATH=$PATH:/opt/mssql-tools/bin